# grimbots.service
[Unit]
Description=GrimBots - Sistema de Bots Telegram
Documentation=https://github.com/grimbots/core
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=simple

# ✅ REALIDADE DO SERVIDOR: /root/grimbots
WorkingDirectory=/root/grimbots
User=root
Group=root

# ✅ EnvironmentFile na localização real
EnvironmentFile=/root/grimbots/.env

# ✅ Criar diretório de logs se não existir
ExecStartPre=/bin/bash -c 'mkdir -p /root/grimbots/logs && chmod 755 /root/grimbots/logs'

# ✅ Gunicorn com logs configurados em /root/grimbots/logs/
ExecStart=/root/grimbots/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 4 \
    --worker-class sync \
    --worker-connections 1000 \
    --timeout 120 \
    --keep-alive 5 \
    --max-requests 1000 \
    --max-requests-jitter 50 \
    --access-logfile /root/grimbots/logs/access.log \
    --error-logfile /root/grimbots/logs/error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance \
    wsgi:app

# ✅ Graceful shutdown
ExecStop=/bin/kill -s TERM $MAINPID
ExecReload=/bin/kill -s HUP $MAINPID

# ✅ Restart automático
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# ✅ Timeout para startup
TimeoutStartSec=60
TimeoutStopSec=30

# ✅ Limites de recursos (ajustar conforme VPS)
LimitNOFILE=65535
LimitNPROC=4096

# ✅ Isolamento de segurança (permissivo para root em /root)
PrivateTmp=no
ProtectSystem=no
ProtectHome=no
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
