# grimbots.service
[Unit]
Description=GrimBots - Telegram Bot Platform
Documentation=https://grimbots.com/docs
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=notify
User=grpay
Group=grpay
WorkingDirectory=/opt/grpay

# Environment File (KEY=VALUE format only - no multiline!)
EnvironmentFile=/opt/grpay/.env

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/grpay /tmp /opt/grpay/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource Limits
LimitAS=2G
LimitNOFILE=65535
LimitNPROC=4096

# Process Management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Health Check (requires app.py health endpoint)
ExecStartPre=/bin/sh -c 'until curl -s localhost:6379 > /dev/null 2>&1 || ! grep -q REDIS_URL /opt/grpay/.env; do sleep 1; done'
ExecStart=/opt/grpay/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 3 \
    --threads 4 \
    --worker-class gthread \
    --worker-connections 1000 \
    --max-requests 10000 \
    --max-requests-jitter 1000 \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --access-logfile /var/log/grpay/access.log \
    --error-logfile /var/log/grpay/error.log \
    --capture-output \
    --enable-stdio-inheritance \
    --preload \
    "app:create_app()"

# Health check endpoint
ExecStartPost=/bin/sh -c 'sleep 5 && curl -s http://localhost:5000/health || echo "Health check failed"'

[Install]
WantedBy=multi-user.target
