# rq-scheduler.service
[Unit]
Description=RQ Scheduler - GrimBots Job Scheduler (Native RQ)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/grimbots
Environment=PYTHONPATH=/root/grimbots
Environment=PATH=/root/grimbots/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=REDIS_URL=redis://localhost:6379/0

# ✅ Comando NATIVO RQ: O scheduler é built-in no RQ 2.7.0+
ExecStart=/root/grimbots/venv/bin/rq scheduler --url redis://localhost:6379/0

# ✅ RESILIÊNCIA MÁXIMA
Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# ✅ LOGS VIA JOURNALD
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rq-scheduler

# ✅ LIMITE DE RECURSOS
MemoryMax=256M
CPUQuota=100%

# ✅ GRACEFUL SHUTDOWN
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
