# deploy/systemd/rq-worker@.service
[Unit]
Description=RQ Worker %I
After=network.target redis.service grimbots.service
Wants=redis.service
PartOf=grimbots.service
Documentation=https://python-rq.org/

[Service]
Type=simple
User=grimbots
Group=grimbots
WorkingDirectory=/opt/grimbots

# Environment
Environment="PATH=/opt/grimbots/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/grimbots"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="DATABASE_URL=postgresql://grimbots:YOUR_PASSWORD@localhost:5432/grimbots"

# RQ Worker command
# %i será substituído pelo nome após @ (ex: tasks-1, gateway-1)
ExecStart=/opt/grimbots/venv/bin/python /opt/grimbots/start_rq_worker.py %i

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Resource limits
LimitNOFILE=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rq-worker-%i

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

